# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:octane]
; Laravel Octane with Swoole - keeps the app in memory for high performance
command=php /var/www/html/artisan octane:start --server=swoole --host=127.0.0.1 --port=8000 --max-requests=500
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=5
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:reverb]
; Laravel Reverb WebSocket Server
command=php /var/www/html/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
; Laravel Queue Worker - Processes jobs from Redis queues
; Disabled: Horizon manages queue workers instead
command=php /var/www/html/artisan queue:work redis --queue=emails,tokens,default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Environment variables for the queue worker
environment=APP_ENV="production",APP_DEBUG="false"

[program:horizon]
; Laravel Horizon - Queue monitor and worker manager
; Manages queue workers and provides a monitoring dashboard at /horizon
command=php /var/www/html/artisan horizon
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Environment variables
environment=APP_ENV="production"
