# docker/supervisor/surprisemoi-queue.conf
; =============================================================================
; SurpriseMoi Queue Worker Supervisor Configuration
; Location: /etc/supervisor/conf.d/surprisemoi-queue.conf
; =============================================================================
;
; This configuration manages Laravel queue workers using Supervisor.
; Queue workers process jobs from Redis queues (emails, tokens, default).
;
; Installation:
; 1. Copy this file to /etc/supervisor/conf.d/surprisemoi-queue.conf
; 2. Run: sudo supervisorctl reread
; 3. Run: sudo supervisorctl update
; 4. Run: sudo supervisorctl start surprisemoi-queue:*
;
; =============================================================================

[program:surprisemoi-queue]
; Command to run the queue worker
command=php /var/www/html/artisan queue:work redis --queue=emails,tokens,default --sleep=3 --tries=3 --max-time=3600

; Directory where the command should be executed
directory=/var/www/html

; User to run the process as (should be www-data or equivalent)
user=www-data

; Process name (will be supervisemoi-queue:surprisemoi-queue_00, surprisemoi-queue:surprisemoi-queue_01, etc.)
process_name=%(program_name)s_%(process_num)02d

; Number of worker processes to spawn
numprocs=2

; Where to redirect STDOUT and STDERR logs
stdout_logfile=/var/log/supervisor/surprisemoi-queue-out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/surprisemoi-queue-err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

; Redirect stderr to stdout (combine logs)
redirect_stderr=true

; Whether to start the process when supervisor starts
autostart=true

; Whether to restart the process after unexpected exit
autorestart=true

; Stop signal to send (TERM allows graceful shutdown)
stopasgroup=true
killasgroup=true

; Exit codes that should be considered "expected" (not a crash)
exitcodes=0

; Stop timeout in seconds (allows graceful shutdown of workers)
stopwaitsecs=30

; Environment variables
environment=APP_ENV="production",APP_DEBUG="false",DB_CONNECTION="pgsql",REDIS_CLIENT="phpredis"

; Priority (lower number = higher priority)
priority=999

; =============================================================================
; Alternative: Redis-based queue with separate workers for each queue
; Uncomment below for dedicated workers per queue
; =============================================================================

; [program:surprisemoi-queue-emails]
; command=php /var/www/html/artisan queue:work redis --queue=emails --sleep=3 --tries=3 --max-time=3600
; directory=/var/www/html
; user=www-data
; process_name=%(program_name)s_%(process_num)02d
; numprocs=1
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; priority=998

; [program:surprisemoi-queue-tokens]
; command=php /var/www/html/artisan queue:work redis --queue=tokens --sleep=3 --tries=3 --max-time=3600
; directory=/var/www/html
; user=www-data
; process_name=%(program_name)s_%(process_num)02d
; numprocs=1
; autostart=true
; autorestart=true
; stopasgroup=true
; killasgroup=true
; priority=997

; =============================================================================
; Supervisor Commands Reference
; =============================================================================
;
; sudo supervisorctl reread           # Reload configuration
; sudo supervisorctl update           # Apply configuration changes
; sudo supervisorctl status           # Check process status
; sudo supervisorctl start surprisemoi-queue:*  # Start all workers
; sudo supervisorctl stop surprisemoi-queue:*   # Stop all workers
; sudo supervisorctl restart surprisemoi-queue:* # Restart all workers
; sudo supervisorctl tail -f surprisemoi-queue_00  # View logs in real-time
;
; =============================================================================
